# systemd/azazel-web.service
[Unit]
Description=Azazel-Edge Web UI Backend (Flask)
After=network-online.target azctl-unified.service
Wants=network-online.target azctl-unified.service

[Service]
Type=simple
WorkingDirectory=/opt/azazel
Environment="AZAZEL_WEB_HOST=127.0.0.1"
Environment="AZAZEL_WEB_PORT=8084"
EnvironmentFile=-/etc/default/azazel-edge
ExecStart=/usr/bin/env python3 /opt/azazel/azazel_web/app.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=azazel-web

# Security hardening
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
