# systemd/azctl-serve.service
[Unit]
Description=Azazel control daemon (serve mode)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/azazel/Azazel-Pi
Environment=PYTHONPATH=/home/azazel/Azazel-Pi
# Ensure decisions directory exists before starting
ExecStartPre=/bin/mkdir -p /var/log/azazel
ExecStart=/usr/bin/env python3 -m azctl.cli serve --suricata-eve /var/log/suricata/eve.json --decisions-log /var/log/azazel/decisions.log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
