# systemd/azctl-unified.service
[Unit]
Description=Azazel Unified Control Daemon with AI Edge Computing
Documentation=https://github.com/01rabbit/Azazel-Pi
After=network-online.target suricata.service azazel-ai-services.service azazel-wan-manager.service
Wants=network-online.target azazel-ai-services.service azazel-wan-manager.service
Requires=suricata.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/azazel/Azazel-Pi
Environment=PYTHONPATH=/home/azazel/Azazel-Pi

# Ensure required directories exist (created as root).
# This service runs as root to allow nft/tc operations; adjust permissions as needed.
ExecStartPre=+/bin/mkdir -p /var/log/azazel

# Unified startup command combining serve mode with config file support
ExecStart=/usr/bin/env python3 -u -m azctl.cli serve --config /etc/azazel/azazel.yaml --suricata-eve /var/log/suricata/eve.json --decisions-log /var/log/azazel/decisions.log

# Service management
# Run as root to allow privileged network operations (nft/tc). Be aware of
# the security implications: consider using a restricted capability drop or
# an additional wrapper if you want to limit privileges later.
# (User/Group intentionally omitted to run as root)
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE CAP_SYS_ADMIN
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
NoNewPrivileges=no
Restart=on-failure
RestartSec=5
KillMode=mixed
TimeoutStopSec=30

# Reduced security settings for initial deployment
# TODO: Enhance security settings after verification

[Install]
WantedBy=multi-user.target
Also=suricata.service vector.service
