# systemd/link-opt.service
[Unit]
Description=Link /opt/azazel/config files to repository
After=network.target

[Service]
Type=oneshot
ExecStart=/home/azazel/Azazel-Edge/scripts/link_opt_to_repo.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
