# systemd/pc-monitor.timer
[Unit]
Description=Run MacBook detector every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
