# systemd/vector.service
[Unit]
Description=Vector log shipper
After=network.target

[Service]
Type=simple
# Try the common install locations: prefer /usr/local/bin (tarball) then /usr/bin (apt)
ExecStart=/bin/sh -c '/usr/local/bin/vector --config /etc/azazel/vector/vector.toml || /usr/bin/vector --config /etc/azazel/vector/vector.toml'
Restart=on-failure

[Install]
WantedBy=azctl.target
